# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisor/supervisord.pid
user=www-data

; All-in-one supervisor config: Octane + Horizon + scheduler in a single
; container. Use supervisord-web.conf / supervisord-worker.conf to split
; roles across multiple containers. Select via the ROLE env var in
; entrypoint.sh (web | worker | all).

[unix_http_server]
file=/var/run/supervisor/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor/supervisor.sock

[program:octane]
command=php /var/www/html/artisan octane:frankenphp --host=0.0.0.0 --port=80 --workers=%(ENV_OCTANE_WORKERS)s --max-requests=%(ENV_OCTANE_MAX_REQUESTS)s
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=10
user=www-data
priority=5

[program:horizon]
command=php /var/www/html/artisan horizon
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=3600
stopsignal=TERM
user=www-data
priority=10

[program:scheduler]
command=sh -c "while true; do php /var/www/html/artisan schedule:run --no-interaction; sleep 60; done"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=www-data
priority=20
